# backend.service
[Unit]
Description=backend

[Service]
User=terraform
WorkingDirectory= /home/terraform/TeamAPoolProjectBackend
ExecStart=/usr/bin/node Index.js


[Install]
WantedBy=multi-user.target
