# backend.service
[Unit]
Description=backend

[Service]
User=terraform
WorkingDirectory=~/TeamAPoolProjectBackend
ExecStart=/usr/bin/node ~/TeamAPoolProjectBackend/Index.js


[Install]
WantedBy=multi-user.target
